# packages/sprixe-image/stage-sprixe/01-kiosk-config/files/sprixe-watchdog.timer
[Unit]
Description=Run Sprixe watchdog every 30s

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s
Unit=sprixe-watchdog.service

[Install]
WantedBy=timers.target
