# packages/sprixe-image/stage-sprixe/01-kiosk-config/files/sprixe-watchdog.service
[Unit]
Description=Sprixe Kiosk Health Watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sprixe-watchdog.sh
